# crontab
@reboot sleep 60 && /home/pi/bojan-tv/download.sh
